# backend/llm-ocr.service
[Unit]
Description=LLM-OCR Flask API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/llm-ocr/backend
Environment=PATH=/path/to/llm-ocr/backend/venv/bin
ExecStart=/path/to/llm-ocr/backend/venv/bin/gunicorn --workers 4 --bind 127.0.0.1:8000 --timeout 120 app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target